# root/etc/uci-defaults/luci-snmpd-tn
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@snmpd[-1]
	add ucitrack snmpd
	set ucitrack.@snmpd[-1].init=snmpd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
